# rc.local
#!/bin/sh

PARTITIONS="`kenv custom.partitions | tr _ \  `"
[ -z "$PARTITIONS" ] && PARTITIONS=DEFAULT
export PARTITIONS
bsdinstall script /etc/installerconfig
